# config/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:httpd]
command=service httpd start

[program:nagios]
command=service nagios start
